# quota-matching-tool.service
[Unit]
Description=企业定额自动套用工具
After=network.target mysql.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/quota-matching-tool
ExecStart=/usr/bin/java -Xms512m -Xmx1024m -XX:+UseG1GC -Djava.awt.headless=true -Dfile.encoding=UTF-8 -Duser.timezone=Asia/Shanghai -jar /opt/quota-matching-tool/target/quota-matching-tool-1.0.0.jar --spring.profiles.active=prod
ExecStop=/bin/kill -15 $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# 环境变量配置（根据需要修改）
Environment="DB_HOST=localhost"
Environment="DB_PORT=3306"
Environment="DB_NAME=quota_db"
Environment="DB_USERNAME=root"
Environment="DB_PASSWORD="
Environment="LOG_PATH=/var/log/quota-matching-tool"

[Install]
WantedBy=multi-user.target

